# deploy/systemd/driftwood-data.service
[Unit]
Description=Driftwood live tide and wave data fetcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/driftwood/scripts
# Stagger startup after boot (replaces old rc.local sleep)
ExecStartPre=/bin/sleep 60
ExecStart=/usr/bin/python3 /home/pi/driftwood/scripts/live_data_stream.py
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
